# dot_config/systemd/user/appimage-manager.path
[Unit]
Description=Monitor for managing AppImages

[Path]
PathChanged=%h/.config/appimages.json

[Install]
WantedBy=default.target
