# dot_config/systemd/user/audible-tool.path
[Unit]
Description=Service for processing Audible AAX files
Requires=audible-tool.service

[Path]
PathChanged=%h/Downloads

[Install]
WantedBy=default.target
