# dot_config/systemd/user/flatpak-manager.path
[Unit]
Description=Monitor for managing Flatpaks

[Path]
PathChanged=%h/.config/flatpaks.json

[Install]
WantedBy=default.target
